# outputs/squashfs-root/usr/lib/opkg/info/adb.list
/bin/adb_shell
/etc/init.d/adbd
/bin/adbd
/etc/adb_profile
