# outputs/squashfs-root/usr/lib/opkg/info/alarm.list
/usr/bin/alarmd
/etc/workday/unnomal_workday.cfg
/etc/workday/unnomal_holiday.cfg
/etc/init.d/alarm
