# outputs/squashfs-root/usr/lib/opkg/info/base-files.list
/etc/openwrt_version
/etc/banner.failsafe
/etc/resolv.conf
/etc/init.d/done
/bin/config_generate
/var
/etc/init.d/system
/etc/preinit
/bin/config_update.sh
/etc/localtime
/etc/shells
/lib/preinit/10_sysinfo
/lib/preinit/70_initramfs_test
/etc/device_info
/etc/protocols
/etc/uci-defaults/12_network-generate-ula
/lib/preinit/10_indicate_preinit
/lib/preinit/79_format_partition
/etc/config/system
/etc/fstab
/etc/sysctl.conf
/sbin/sysupgrade
/lib/preinit/99_10_run_init
/sbin/led.sh
/lib/preinit/50_indicate_regular_preinit
/etc/banner
/etc/shadow
/etc/rc.common
/etc/uci-defaults/10_migrate-shadow
/etc/board.d/99-default_network
/usr/share/mico/version
/lib/preinit/10_indicate_failsafe
/init
/etc/board.d/02_network
/etc/sysctl.d/local.conf
/lib/preinit/80_mount_root
/etc/rc.local
/lib/functions/uci-defaults.sh
/etc/init.d/gpio_switch
/lib/functions/network.sh
/etc/init.d/sysctl
/sbin/hotplug-call
/etc/crontabs/root
/etc/config/fstab
/sbin/imiflash
/etc/asound.conf
/etc/rc.button/failsafe
/etc/services
/lib/functions/system.sh
/bin/board_detect
/lib/upgrade/common.sh
/etc/rc.button/power
/etc/sysupgrade.conf
/etc/rc.button/rfkill
/lib/upgrade/keep.d/base-files
/lib/preinit/40_run_failsafe_hook
/etc/group
/etc/rc.button/reset
/lib/preinit/30_failsafe_wait
/sbin/firstboot
/bin/login.sh
/etc/hotplug.d/net/00-sysctl
/etc/hosts
/bin/getmac.sh
/lib/functions.sh
/etc/iproute2/rt_tables
/etc/uci-defaults/11_migrate-sysctl
/lib/functions/leds.sh
/rom/note
/etc/inittab
/etc/openwrt_release
/bin/bt_connect_json.sh
/sbin/wifi
/etc/opkg/keys/af22f7a88858c8e9
/etc/profile
/lib/preinit/99_10_failsafe_login
/etc/init.d/sysfixtime
/etc/TZ
/etc/init.d/umount
/lib/functions/service.sh
/lib/upgrade/keep.d/base-files-essential
/etc/passwd
/lib/preinit/02_default_set_state
/etc/mtab
/etc/init.d/boot
/etc/init.d/check_mac
/lib/functions/preinit.sh
/bin/ipcalc.sh
/etc/config/bluetooth
/etc/diag.sh
