# outputs/squashfs-root/usr/lib/opkg/info/busybox.list
/bin/dmesg
/usr/bin/id
/bin/login
/bin/rm
/sbin/ip
/bin/sh
/bin/su
/usr/bin/test
/bin/zcat
/bin/fgrep
/bin/ash
/sbin/halt
/sbin/reboot
/bin/true
/sbin/route
/usr/bin/md5sum
/usr/bin/free
/usr/sbin/ntpd
/usr/bin/time
/usr/bin/passwd
/usr/bin/traceroute
/bin/vi
/bin/pwd
/usr/bin/nc
/bin/chown
/etc/init.d/cron
/bin/mount
/bin/sync
/usr/bin/dirname
/sbin/hwclock
/sbin/udhcpc
/usr/bin/traceroute6
/bin/ping6
/usr/bin/basename
/bin/echo
/usr/bin/less
/usr/sbin/chroot
/usr/sbin/crond
/usr/bin/head
/bin/gunzip
/bin/netstat
/bin/busybox
/bin/grep
/usr/bin/clear
/usr/bin/seq
/bin/sleep
/bin/nice
/bin/kill
/bin/mktemp
/usr/bin/tr
/bin/mknod
/bin/rmdir
/usr/bin/uniq
/usr/bin/tee
/usr/bin/bzcat
/usr/bin/hexdump
/etc/init.d/sysntpd
/bin/date
/usr/bin/wc
/usr/bin/[
/usr/bin/nslookup
/usr/bin/pgrep
/bin/umount
/bin/netmsg
/usr/bin/tail
/usr/sbin/ntpd-hotplug
/usr/bin/mkfifo
/bin/chgrp
/bin/cp
/sbin/poweroff
/usr/bin/[[
/bin/dd
/bin/lock
/bin/df
/bin/false
/usr/bin/killall
/usr/bin/sort
/bin/mkdir
/sbin/ifconfig
/usr/bin/cmp
/usr/bin/reset
/sbin/start-stop-daemon
/usr/sbin/brctl
/bin/cat
/bin/touch
/sbin/mkswap
/bin/tar
/usr/bin/yes
/usr/bin/top
/sbin/switch_root
/usr/bin/expr
/usr/bin/env
/sbin/sysctl
/bin/sed
/usr/bin/logger
/bin/chmod
/usr/bin/awk
/bin/uname
/usr/bin/which
/usr/bin/find
/bin/ln
/bin/ls
/sbin/pivot_root
/usr/bin/cut
/usr/bin/du
/bin/mv
/usr/bin/crontab
/usr/bin/printf
/bin/ping
/usr/bin/timeout
/usr/bin/bunzip2
/usr/bin/uptime
/bin/fsync
/bin/pidof
/bin/egrep
/usr/bin/readlink
/usr/bin/strings
/bin/ps
/usr/bin/xargs
/bin/gzip
