# outputs/squashfs-root/usr/lib/opkg/info/dnsmasq.list
/etc/dnsmasq.conf
/etc/hotplug.d/iface/25-dnsmasq
/etc/config/dhcp
/etc/init.d/dnsmasq
/usr/sbin/dnsmasq
