# outputs/squashfs-root/usr/lib/opkg/info/mediaplayer.list
/usr/lib/libxiaomimediaplayer.so
/etc/init.d/mediaplayer.init
/usr/lib/libffmpeg-miplayer.so
/usr/bin/mphelper
/usr/bin/mediaplayer
/usr/bin/check_mediaplayer_status
