# outputs/squashfs-root/usr/lib/opkg/info/messagingagent.list
/usr/bin/matool_get_id_for_vendor
/usr/bin/matool_upload_log
/usr/bin/matool_get_rom_version
/usr/share/mico/messaging/serverkey_2.pub
/usr/bin/matool_get_device_status
/usr/bin/matool_request_access_token
/usr/bin/matool_set_device_status
/usr/bin/matool_get_ai_service_token
/usr/lib/libxiaomi_mico.so
/usr/lib/libxiaomi_mosquitto.so
/usr/lib/libxiaomi_http.so
/usr/bin/matool_time_sync
/usr/bin/matool_notify_app
/usr/bin/matool_upload_action_log
/usr/bin/matool
/usr/bin/matool_get_miot_did
/usr/bin/quark
/usr/bin/matool_get_round_trip_delay_time
/usr/bin/matool_get_rom_channel
/usr/bin/matool_get_album_song
/usr/bin/matool_get_station_sound_list
/usr/bin/matool_get_mac
/usr/lib/libquark_lib.so
/usr/bin/matool_localtime_r
/usr/bin/matool_upload_miot_token
/usr/bin/matool_get_hardware
/etc/init.d/messagingagent
/usr/lib/liblibma.so
/usr/bin/statpoints_daemon
/usr/bin/matool_get_playlist_song
/usr/bin/matool_generate_uuid
/usr/lib/libxiaomi_crypto.so
/usr/lib/libxiaomi_json.so
/usr/bin/matool_check_registration_status
/usr/bin/matool_get_channel_secret
/usr/sbin/statpoints_clear.sh
/usr/bin/matool_get_device_playlists
/usr/bin/matool_get_music_recommendation
/usr/bin/matool_get_sn
/usr/share/mico/messaging/serverkey.pub
/usr/bin/matool_get_broker_address
/etc/init.d/statpoints_daemon
/etc/quark.cfg
/usr/bin/quarkd
/usr/bin/matool_get_playlistv2_song
/usr/lib/libxiaomi_utils.so
/usr/bin/matool_get_device_id
/usr/bin/matool_check_upgrade
/usr/share/mico/messaging/messaging
/usr/bin/messagingagent
/usr/share/mico/messaging/messaging.conf
/usr/share/mico/messaging/actionkey.pub
/usr/bin/matool_compute_md5
/usr/bin/statpoints
/usr/bin/matool_get_super_admin
/usr/bin/matool_collect_stat_points
/usr/bin/matool_music_action_log
/usr/bin/matool_get_tencent_account
