# outputs/squashfs-root/usr/lib/opkg/info/mibrain-service.list
/etc/init.d/mibrain_service
/usr/lib/libmibrain-util.so
/usr/bin/mibrain_service
