# outputs/squashfs-root/usr/lib/opkg/info/mibt.list
/etc/bt/mibt_uuid.json
/usr/bin/mibt_ble
/usr/bin/app_ble
/etc/init.d/bluetooth
/usr/bin/app_avk
/usr/bin/app_manager
/etc/bt/mibt_config.json
/usr/bin/bsa_server
