# outputs/squashfs-root/usr/lib/opkg/info/mico-misc.list
/bin/ota
/usr/share/sound/mic_off.wav
/usr/share/sound/bootup.mp3
/bin/bind_device.sh
/bin/show_led
/usr/share/sound/shutdown.mp3
/usr/share/sound/wakeup_mitv.mp3
/bin/ntpsetclock
/usr/share/sound/mic_on.wav
/bin/wpa_action.sh
/usr/bin/reboot.sh
/usr/share/sound/internet_disconnect.mp3
/lib/mico_common.sh
/usr/sbin/syslog_rotate
/usr/sbin/easy_logcut
/usr/sbin/switch_channel
/usr/share/sound/network_done.wav
/bin/wl
/usr/share/sound/mibrain_service_timeout.mp3
/usr/share/sound/reset.wav
/usr/share/sound/mibrain_network_unreachable.mp3
/usr/sbin/supervisord
/usr/share/sound/wakeup_wozai.wav
/usr/share/sound/bluetooth_failure.wav
/usr/share/sound/brainTurnUpVol.wav
/usr/share/sound/weak_network.mp3
/usr/share/sound/alarmDefault.mp3
/usr/share/sound/upgrade_now.mp3
/usr/share/sound/wakeup_ei.wav
/usr/share/sound/wakeup2.wav
/usr/share/sound/welcome.opus
/etc/init.d/boot_check
/usr/share/sound/timer_default.opus
/bin/easy_play
/usr/share/sound/enter_config_mode.mp3
/usr/share/sound/wifi_disconnect.mp3
/usr/share/sound/notice.wav
/usr/bin/mibrain_level
/usr/share/sound/mibrain_service_unreachable.mp3
/etc/init.d/wireless
/usr/sbin/collect_log.sh
/usr/share/sound/reset_wait.mp3
/bin/simple_dhcp.sh
/usr/share/sound/bluetooth_disconnect.wav
/etc/init.d/led
/usr/share/sound/bluetooth_connect.wav
/bin/wifi_check_ccmp.sh
/usr/share/sound/upgrade_later.mp3
/usr/share/sound/setup_failure.wav
/etc/wifi/wpa_supplicant_backup.conf
/bin/shut_led
/usr/bin/ssh_enable
/usr/share/sound/mibrain_connect_timeout.mp3
/etc/asound.state
/usr/sbin/syslog_sync.sh
/etc/hotplug.d/iface/40-ntpd
/usr/share/sound/key_prev_next.wav
/usr/share/sound/unknown_domain.mp3
/etc/init.d/alsa
/usr/bin/experience_plan
/usr/share/sound/no_channel.opus
/usr/share/sound/volume.wav
/usr/share/sound/brainTurnDownVol.wav
/bin/notify.sh
/usr/sbin/network_probe.sh
/usr/share/sound/wakeup_zai.wav
/usr/bin/mibrain_net_check
/usr/share/sound/first_voice.wav
/usr/share/sound/bluetooth_noPhone.wav
/usr/share/sound/unknown_action.mp3
/etc/init.d/touchpad
/usr/share/sound/service_timeout.mp3
/usr/share/sound/network_done_init.wav
/etc/ssl/certs/ca-cert.crt
/bin/set_upgrade_status
/etc/wifi/hostapd.conf
/usr/share/sound/wakeup.wav
/usr/bin/gen-tts
/bin/touchpad
/bin/unbind_device.sh
/etc/init.d/dhcpc
/usr/bin/config_mode
/bin/ledserver
