# outputs/squashfs-root/usr/lib/opkg/info/mipns-general.list
/etc/init.d/pns
/bin/tplay
/bin/wuw_upload.sh
/bin/wakeup.sh
