# outputs/squashfs-root/usr/lib/opkg/info/mipns-helper.list
/usr/bin/pns_ubus_helper
/etc/init.d/pns_ubus_helper
/usr/bin/pns_upload_helper
