# outputs/squashfs-root/usr/lib/opkg/info/misc-upgrade.list
/sbin/write_private
/sbin/write_misc
/sbin/aw_upgrade_vendor_default.sh
/sbin/aw_upgrade_process.sh
/sbin/aw_reboot.sh
/etc/init.d/aw_upgrade_autorun
/sbin/aw_upgrade_image.sh
/sbin/aw_upgrade_normal.sh
/sbin/read_misc
/etc/init.d/aw_fstab
/sbin/aw_upgrade_utils.sh
/sbin/aw_upgrade_log.sh
/sbin/read_private
/sbin/resetsn.sh
