# outputs/squashfs-root/usr/lib/opkg/info/nano-httpd.list
/etc/init.d/nano_httpd
/usr/bin/nano_httpd
