# outputs/squashfs-root/usr/lib/opkg/info/netifd.list
/lib/netifd/proto/dhcp.sh
/etc/init.d/network
/lib/netifd/netifd-proto.sh
/etc/hotplug.d/iface/00-netstate
/sbin/ifdown
/lib/netifd/dhcp.script
/usr/share/udhcpc/default.script
/lib/network/config.sh
/sbin/netifd
/sbin/devstatus
/lib/netifd/utils.sh
/lib/netifd/netifd-wireless.sh
/sbin/ifstatus
/sbin/ifup
