# outputs/squashfs-root/usr/lib/opkg/info/quickplayer.list
/usr/bin/quickplayer
/etc/init.d/quickplayer
