# outputs/squashfs-root/usr/lib/opkg/info/syslog-ng.list
/etc/syslog-ng.conf
/etc/init.d/syslog-ng
/usr/sbin/syslog-ng
