# outputs/squashfs-root/usr/lib/opkg/info/upnp-disc.list
/etc/init.d/mitv-disc
/usr/bin/upnp-disc
